# templates/eve/etc/systemd/system/docker-update-containers.service
[Unit]
Description=Update all Docker containers
OnFailure=unit-fail-mail@%n.service

[Service]
WorkingDirectory={{ ansible_user_dir }}/{{ ansible_hostname }}
ExecStartPre={{ docker_compose_path }} pull
ExecStart={{ docker_compose_path }} up -d --build
ExecStartPost=/usr/bin/docker image prune -f
ExecStartPost=/usr/bin/docker builder prune -f