# templates/eve/etc/systemd/system/nextcloud-cron.service
[Unit]
Description=Nextcloud cron.php job

[Service]
ExecStart=/usr/bin/docker exec --user www-data nextcloud php -f cron.php
