# templates/eve/etc/systemd/system/nextcloud-preview-generator.service
[Unit]
Description=Nextcloud generate previews

[Service]
ExecStart=/usr/bin/docker exec --user www-data nextcloud php occ preview:pre-generate
