# server/templates/fiona/etc/systemd/system/update-ip-on-dynv6.timer
[Unit]
Description=Update IP on DynV6

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
