# server/templates/kubo/etc/systemd/system/snmpd.service.d/10-restart_on_failure.conf
[Unit]
RequiresMountsFor=/data/bulk /data/important
