# server/templates/hosts/fiona-code-runner/etc/systemd/system/dind-prune.timer
[Unit]
Description=Prune Docker-in-Docker

[Timer]
OnCalendar=22:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
