# server/roles/jenswbe.powersavings/files/autoaspm.service
[Unit]
Description=AutoASPM

[Service]
Type=oneshot
ExecStart=/usr/sbin/autoaspm.py

[Install]
WantedBy=multi-user.target
