# server/templates/hosts/kubo-media/etc/systemd/system/docker.service.d/10-wait_on_data_mounts.conf
[Unit]
Wants=data-bulk.mount
After=data-bulk.mount

Wants=data-important.mount
After=data-important.mount
