# server/templates/docker-update-containers.timer
[Unit]
Description=Update Docker containers every night

[Timer]
OnCalendar={{ var_time }}
Persistent=false

[Install]
WantedBy=timers.target
