# server/templates/eve/etc/systemd/system/nextcloud-preview-generator.timer
[Unit]
Description=Nextcloud generate previews every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
