# templates/kubo/etc/systemd/system/nextcloud-update-apps.timer
[Unit]
Description=Update Nextcloud apps every night

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=false

[Install]
WantedBy=timers.target
