# server/templates/hosts/eve/etc/systemd/system/docker-update-containers.timer
[Unit]
Description=Update Docker containers every night

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=false

[Install]
WantedBy=timers.target
