# server/templates/hosts/eve/etc/systemd/system/miniflux-reset-feed-errors.timer
[Unit]
Description=Reset MiniFlux feed errors

[Timer]
OnCalendar=Mon 01:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
