# server/templates/hosts/eve/etc/systemd/system/miniflux-reset-feed-errors.service
[Unit]
Description=Reset MiniFlux feed errors
OnFailure=unit-fail-mail@%n.service

[Service]
Type=exec
ExecStart={{ docker_path }} exec miniflux miniflux --debug --reset-feed-errors
