# server/templates/hosts/eve/etc/systemd/system/nextcloud-db-add-missing-indices.timer
[Unit]
Description=Add missing indexes in DB for Nextcloud every night

[Timer]
OnCalendar=01:00
Persistent=false

[Install]
WantedBy=timers.target
