# server/templates/hosts/eve/etc/systemd/system/nextcloud-calcardbackup.timer
[Unit]
Description=Backup all Nextcloud calendars and contacts at 02h00

[Timer]
OnCalendar=02:15
Persistent=true

[Install]
WantedBy=timers.target
