# server/templates/hosts/eve/etc/systemd/system/vaultwarden-dump-db.timer
[Unit]
Description=Dump Vaultwarden DB every night for backup

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
