# server/templates/hosts/eve/etc/systemd/system/wikijs-dump-db.timer
[Unit]
Description=Dump Wiki.js DB every night for backup

[Timer]
OnCalendar=*-*-* 02:05:00
Persistent=true

[Install]
WantedBy=timers.target
