# desktop/hercules/files/clone-displays.service
[Unit]
Description=Clone virtual and real display

[Service]
Type=exec
ExecStart=/usr/local/bin/clone-displays.sh

[Install]
WantedBy=default.target
