# server/templates/kubo/etc/systemd/system/mnt-btrfs.mount
[Mount]
What=/dev/disk/by-uuid/c44fbbe8-f708-4bbd-9dfc-8a0fa15fc956
Where=/mnt/btrfs
Type=btrfs

[Install]
WantedBy=multi-user.target