# server/templates/hosts/kubo/etc/systemd/system/glitchtip-dump-db.timer
[Unit]
Description=Dump GlitchTip every night for backup

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
