# radio/physical_and_tts/radio.service
[Unit]
Description=Radio
 
[Install]
WantedBy=multi-user.target
 
[Service]
Type=simple
PrivateTmp=true
ExecStart=/usr/bin/radio 
